# rootfs/etc/services.d/10-confd/run
#!/bin/sh

with-contenv

exec /opt/wodby/tools/bin/confd -watch
