# rootfs/etc/services.d/00-dnsmasq/run
#!/bin/sh

if [ -f ./setup ];then
    ./setup
fi
if [ $? = "0" ];then
    exec 2>&1
    exec /opt/wodby/tools/bin/go-dnsmasq --default-resolver --append-search-domains
else
    exit 1
fi